# disk/udev/init.d/udev
#! /bin/bash
#
# random	init script to setup /udev
#
# chkconfig: 2345 20 80
# description: manage user-space device nodes in /udev
#
# 2003-12-09 Eric Sandall <eric@sandall.us>
#  Modified for Source Mage GNU/Linux and the Simpleinit-MSB init system
#
# 2003-12-18 Arwed v. Merkatz <v.merkatz@gmx.net>
#  Removed usage of echo/cut to set DEVPATH
#  Adds all device classes now, not only block and tty


RUNLEVEL=S
DEPENDS="+local_fs hotplug"
PROGRAM=/sbin/udev
SYSFS=/sys

. /etc/init.d/smgl_init

start()
{
  if [ ! -d ${UDEV} ]; then
    mkdir ${UDEV}
  fi
  if [ ! -d ${SYSFS} ]; then
    mkdir ${SYSFS}
  fi
  # mount sysfs, ignore errors (if it already is mounted)
  mount -t sysfs none ${SYSFS} > /dev/null 2>&1
  [ -e ${SYSFS}/devices ] || exit 1
  # propogate /udev from /sys - we only need this while we do not
  # have initramfs and an early user-space with which to do early
  # device bring up
  echo "Creating initial udev device nodes: "
  export ACTION=add
  # add block devices and their partitions
  for i in ${SYSFS}/block/*; do
    # add each drive
    export DEVPATH=${i#${SYSFS}}
    $PROGRAM block &

    # add each partition, on each device
    for j in $i/*; do
      if [ -f $j/dev ]; then
        export DEVPATH=${j#${SYSFS}}
        $PROGRAM block &
      fi
    done
  done
  # all other device classes
  for i in ${SYSFS}/class/*; do
    for j in $i/*; do
      if [ -f $j/dev ]; then
        export DEVPATH=${j#${SYSFS}}
        CLASS=${i#${SYSFS}/class/}
        $PROGRAM $CLASS &
      fi
    done
  done
}

stop()
{
  # be careful
  if [ ${UDEV} -a "${UDEV}" != "/" ]; then
    # clear out /udev
    rm -rf ${UDEV}/*
  fi
}

status()
{
  if [ -d ${UDEV} ]; then
    echo "the udev device node directory exists"
  else
    echo "the udev device node directory does not exist"
  fi
}
